# home/.config/bash/profile.d/50-go.sh
# shellcheck shell=sh

export GOPATH=$HOME/.local/lib/go
add_to_path "$GOPATH/bin"
